# deploy/systemd/drone-slam-vio-nav.service
# deploy/systemd/drone-slam-vio-nav.service
# Process 3 — SLAM / VIO / Navigation: visual-inertial odometry, pose estimation.
#
# Depends on perception for fused detections.
# Relationship: slam_vio_nav requires perception (Requires/After), but does not restart perception if it dies.

[Unit]
Description=Drone SLAM/VIO/Nav (P3)
Documentation=https://github.com/nmohamaya/companion_software_stack
PartOf=drone-stack.target
After=network-online.target drone-perception.service
Requires=drone-perception.service

[Service]
Type=simple
ExecStart=/opt/drone/bin/slam_vio_nav --config /etc/drone/config.json --json-logs
Restart=on-failure
RestartSec=2s

# Resource limits
MemoryMax=1G
CPUQuota=200%

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/dev/shm /tmp /var/log/drone
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=drone-slam-vio-nav

[Install]
WantedBy=drone-stack.target
