# deploy/systemd/drone-perception.service
# deploy/systemd/drone-perception.service
# Process 2 — Perception: object detection, tracking, sensor fusion.
#
# Depends on video_capture for camera frames.
# Restart behavior is governed by this unit's Restart= policy and drone-stack.target.

[Unit]
Description=Drone Perception (P2)
Documentation=https://github.com/nmohamaya/companion_software_stack
PartOf=drone-stack.target
After=network-online.target drone-video-capture.service
Requires=drone-video-capture.service

[Service]
Type=simple
ExecStart=/opt/drone/bin/perception --config /etc/drone/config.json --json-logs
Restart=on-failure
RestartSec=2s

# Resource limits
MemoryMax=1G
CPUQuota=200%

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/dev/shm /tmp /var/log/drone
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=drone-perception

[Install]
WantedBy=drone-stack.target
