# deploy/systemd/drone-mission-planner.service
# deploy/systemd/drone-mission-planner.service
# Process 4 — Mission Planner: waypoint tracking, contingency logic.
#
# Depends on comms (FC link) and slam_vio_nav (pose).
# PartOf=drone-comms.service → cascade restart when comms dies.

[Unit]
Description=Drone Mission Planner (P4)
Documentation=https://github.com/nmohamaya/companion_software_stack
PartOf=drone-stack.target drone-comms.service
After=network-online.target drone-comms.service drone-slam-vio-nav.service
Requires=drone-comms.service drone-slam-vio-nav.service

[Service]
Type=simple
ExecStart=/opt/drone/bin/mission_planner --config /etc/drone/config.json --json-logs
Restart=on-failure
RestartSec=2s

# Resource limits
MemoryMax=256M
CPUQuota=100%

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/dev/shm /tmp /var/log/drone
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=drone-mission-planner

[Install]
WantedBy=drone-stack.target
